# thirdparty/xtide-2.15.1/scripts/systemd/xttpd.socket
[Socket]
ListenStream=80
Accept=false

[Install]
WantedBy=sockets.target
